# deploy/frpc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibeflow.frpc</string>
    <key>ProgramArguments</key>
    <array>
        <string>FRPC_PATH</string>
        <string>-c</string>
        <string>FRPC_TOML_PATH</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>LOG_DIR/frpc.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/frpc.err</string>
</dict>
</plist>
